# .crontab
MAILTO=gurdiga@gmail.com
SHELL=/usr/local/bin/bash
PATH=/usr/local/bin:/bin:/usr/bin:/usr/local/opt/coreutils/libexec/gnubin

BIN_DIRS=/bin /usr/bin /usr/local/bin /usr/local/opt/coreutils/libexec/gnubin
HELP_DIR=/Users/vlad/help

30 * * * * sudo /usr/bin/sntp -sS europe.pool.ntp.org 2>&1 >/dev/null | ifne mail -s 'NTP sync error?' gurdiga@gmail.com
31 10 * * 0-4 { grep -Po '^MANPATH\s+\K.+' /private/etc/man.conf; echo -e ${MANPATH//:/\\n}; } | grep -v X11 | xargs -I{} find {} -printf '\%f\n' | sed -E 's/\.([0-9]\.?[a-z]*)$//' | shuf -n 1 | xargs man | col -b | mail -s 'Daily man' gurdiga@gmail.com
32 10 * * 0-4 cat $HELP_DIR/$(ls -1 $HELP_DIR | head -$((${RANDOM} \% `ls -1 $HELP_DIR | wc -l` + 1)) | tail -1) | mail -s 'Daily builtin' gurdiga@gmail.com
